# henryic/deployment/usbdisk.rules
SUBSYSTEM=="block", SUBSYSTEMS=="usb", ACTION=="add" , RUN+="/home/henryic/git/sandbox/henryic/detect_disk.sh"
